# rjjk_web
#!/bin/bash
#
# rjjk_web      Startup script for the RJJK web application
#
# chkconfig: 345 88 12
# description: RJJK web application
# pidfile: /u/apps/$app/shared/log/mongrel.pid

. /etc/rc.d/init.d/functions

app=rjjk_web
user=capistrano
port=3005
pidfile=${PIDFILE-/u/apps/$app/shared/log/mongrel.pid}

# See how we were called.
case "$1" in
  start)
        su - $user -c "cd /u/apps/$app/current;RAILS_ENV=production ruby script/ferret_start"
        su - $user -c "mongrel_rails start -p $port -e production -c /u/apps/$app/current -d 1>/u/apps/$app/shared/log/stdout.log 2>/u/apps/$app/shared/log/stderr.log"
        ;;
  stop)
        su - $user -c "mongrel_rails stop -c /u/apps/$app/current"
        su - $user -c "cd /u/apps/$app/current;RAILS_ENV=production ruby script/ferret_stop"
        rm -f ${pidfile}
        ;;
  restart)
        su - $user -c "cd /u/apps/$app/current;RAILS_ENV=production ruby script/ferret_stop"
        su - $user -c "cd /u/apps/$app/current;RAILS_ENV=production ruby script/ferret_start"
        su - $user -c "mongrel_rails restart -c /u/apps/$app/current"
        ;;
  status)
        ps -ef | grep $app
        ;;
  *)
        echo $"Usage: $prog {start|stop|restart|status}"
        exit 1
esac

exit 0
